# obd-logger/obd-logger.service
[Unit]
Description=OBD Logger
After=network.target

[Service]
ExecStartPre=/sbin/ip link set can0 up type can bitrate 500000
ExecStart=/home/canbus/Documents/Code/SierraECMLogger/venv/bin/python /home/canbus/Documents/Code/SierraECMLogger/main.py
WorkingDirectory=/home/canbus/Documents/Code/SierraECMLogger
Restart=on-failure
RestartSec=5
User=pi

[Install]
WantedBy=multi-user.target